# packaging/howdy-island.service
[Unit]
Description=howdy-island — face-scan OSD for Howdy
Documentation=https://github.com/sijuz/howdy-island
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/howdy-island
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
